# com.user.cfst.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 任务标识 -->
    <key>Label</key>
    <string>com.user.cfst.update</string>
    
    <!-- 执行程序 -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- ⚠️ 修改: 改为你的实际路径 -->
        <string>/path/to/cfst/auto_update_hosts.sh</string>
    </array>
    
    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <!-- ⚠️ 修改: 改为你的实际目录 -->
    <string>/path/to/cfst</string>
    
    <!-- 每 6 小时执行一次 (6*60*60=21600秒) -->
    <key>StartInterval</key>
    <integer>21600</integer>
    
    <!-- 或者设置为每天固定时间点执行（与 StartInterval 二选一） -->
    <!-- <key>StartCalendarInterval</key> -->
    <!-- <dict> -->
    <!--     <key>Hour</key> -->
    <!--     <integer>3</integer> -->
    <!--     <key>Minute</key> -->
    <!--     <integer>0</integer> -->
    <!-- </dict> -->
    
    <!-- 标准输出日志 -->
    <key>StandardOutPath</key>
    <!-- ⚠️ 修改: 改为你的实际路径 -->
    <string>/path/to/cfst/update.log</string>
    
    <!-- 错误日志 -->
    <key>StandardErrorPath</key>
    <!-- ⚠️ 修改: 改为你的实际路径 -->
    <string>/path/to/cfst/update.error.log</string>
    
    <!-- 如果任务失败，不重试 -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- 网络可达时执行 -->
    <key>NetworkState</key>
    <true/>
</dict>
</plist>
